# backup/openwrt/crontab
# [维护] 每季度首日 (1/4/7/10月1日) 04:30 自动重启，预留 5 秒系统缓冲
# 30 4 1 */3 * sleep 5 && touch /etc/banner && reboot

# [监控] 每小时检查系统状态：存储空间或内存不足时触发邮件告警
0 * * * * /root/sys_monitor.sh

# [计划] 每天 22:00 定时熄灭面板指示灯 (进入睡眠模式)
0 22 * * * echo 0 > /sys/class/leds/blue:status/brightness

# [计划] 每天 07:00 定时开启面板指示灯 (恢复工作状态)
0 7 * * * echo 255 > /sys/class/leds/blue:status/brightness
